# infra/systemd/dms.service
[Unit]
Description=DMS Django via Gunicorn
After=network.target

[Service]
User=dms
Group=dms
WorkingDirectory=/srv/dms/app/backend
EnvironmentFile=/etc/dms.env
ExecStart=/srv/dms/app/.venv/bin/gunicorn backend.wsgi:application \
  --bind 127.0.0.1:8000 \
  --workers 1 \
  --worker-class gthread \
  --threads 8 \
  --timeout 600 \
  --graceful-timeout 30 \
  --max-requests 500 \
  --max-requests-jitter 50

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
