# infra/gunicorn.service
[Unit]
Description=DMS Django via Gunicorn
After=network.target

[Service]
User=dms
Group=dms
WorkingDirectory=/srv/dms/app/backend
EnvironmentFile=/etc/dms.env
ExecStart=/srv/dms/app/.venv/bin/gunicorn backend.wsgi:application \
  --timeout 600 --workers 3 --bind 127.0.0.1:8000 \
  --timeout 300 --graceful-timeout 30 --max-requests 500 --max-requests-jitter 50
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
